# api/dillojs-api.conf
start on filesystem and started networking
respawn
chdir /opt/api
env NODE_ENV=production
env PORT=8080
exec node server.js
